# work/cron/crontab
*/10 * * * * /home/jovyan/work/cron/run.sh ./sample_note.ipynb 2>&1

